# supervisord-nginx.conf
[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
directory=/var/lib/puppetmaster

